# rootfs/etc/services.d/elasticsearch/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid ash \
    /usr/share/elasticsearch/bin/elasticsearch "$@"
